# ansible/roles/appserver/templates/app.service.j2
[Unit]
Description=Simple Java HTTP Server
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/app
ExecStart=/usr/bin/java SimpleHttpServer
Environment=APP_PORT={{ app_port | default('8080') }}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target